# apache-oidc/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Configure and launch Apache
# ==============================================================================
set -e

bashio::log.info "Configuring Apache..."

# Basic Apache config
HA_PORT=$(bashio::core.port)
DOMAIN=$(bashio::config 'domain')
CERT_FILE=$(bashio::config 'ssl_cert_file')
KEY_FILE=$(bashio::config 'ssl_key_file')
HSTS=$(bashio::config 'hsts_header')

# OIDC config
OIDC_METADATA_URL=$(bashio::config 'oidc_metadata_url')
OIDC_CLIENT_ID=$(bashio::config 'oidc_client_id')
OIDC_CLIENT_SECRET=$(bashio::config 'oidc_client_secret')
OIDC_SCOPE=$(bashio::config 'oidc_scope')
OIDC_REDIRECT_URI=$(bashio::config 'oidc_redirect_uri')
OIDC_SESSION_INACTIVITY_TIMEOUT=$(bashio::config 'oidc_inactivity_timeout')

# OIDC generate passphrase
OIDC_CRYPTO_PASSPHRASE=$(bash -c '< /dev/urandom tr -dc _A-Z-a-z-0-9 | head -c64')

# Update Apache site config
sed -i "s#%%CERT_FILE%%#$CERT_FILE#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%KEY_FILE%%#$KEY_FILE#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%DOMAIN%%#$DOMAIN#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%HA_PORT%%#$HA_PORT#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_METADATA_URL%%#$OIDC_METADATA_URL#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_CLIENT_ID%%#$OIDC_CLIENT_ID#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_CLIENT_SECRET%%#$OIDC_CLIENT_SECRET#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_SCOPE%%#$OIDC_SCOPE#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_REDIRECT_URI%%#$OIDC_REDIRECT_URI#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_SESSION_INACTIVITY_TIMEOUT%%#$OIDC_SESSION_INACTIVITY_TIMEOUT#g" /etc/apache2/conf-enabled/home-assistant.conf
sed -i "s#%%OIDC_CRYPTO_PASSPHRASE%%#$OIDC_CRYPTO_PASSPHRASE#g" /etc/apache2/conf-enabled/home-assistant.conf

[ -n "$HSTS" ] && HSTS="Header always set Strict-Transport-Security \"$HSTS\""
sed -i "s#%%HSTS%%#$HSTS#g" /etc/apache2/conf-enabled/home-assistant.conf

# Launch httpd
bashio::log.info "Configuration succeeded. Launching Apache..."
exec /usr/sbin/apache2ctl -DFOREGROUND
